# ansible/roles/common/files/seed-environment.target
[Unit]
Description=/etc/environment available and hostname set
Requires=seed-environment.service
